# systemd/user/arch-dotfiles-sync.timer
[Unit]
Description=Daily Arch dotfile snapshot

[Timer]
OnCalendar=*-*-* 18:00:00
Persistent=true
RandomizedDelaySec=30m
AccuracySec=5m

[Install]
WantedBy=timers.target
